# systemd/systemd-firewalld.service
[Unit]
Description=Systemd Firewall Manager

[Service]
Type=oneshot
ExecStart=/usr/share/systemd-firewald/systemd-firewald.sh
